# etc/service/dnsmasq/run
#!/bin/bash -e

# If internal DNS is provided by Docker, we do not run.
if grep -q 127.0.0.11 /etc/resolv.conf ; then
  # By touching a file in memory, we assure that next time this container is run,
  # it is not present. We have a symlink from /etc/service/dnsmasq/down which points to it.
  # When down file is present (through symlink), this service is not run.
  touch /dev/shm/dnsmasq-disabled
  exit 0
else
  # Has to run as root otherwise it tries to keep some capabilities and it fails inside Docker.
  exec /usr/sbin/dnsmasq --listen-address=127.0.0.11 --addn-hosts=/etc/hosts.extra --keep-in-foreground --bind-interfaces --log-facility=- --user=root 2>&1
fi
