# etc/service/dockergen/run
#!/bin/bash -e

exec /dockergen/docker-gen -config /dockergen/dockergen.conf 2>&1
