# etc/service/dockergen/run
#!/bin/bash -e

exec /dockergen/docker-gen -watch -only-exposed -notify "/usr/sbin/nginx -s reload" /dockergen/nginx.tmpl /etc/nginx/sites-enabled/virtual 2>&1

