# etc/cron.daily/letsencrypt
#!/bin/bash -e

/letsencrypt/dockergen.sh
